# tools/hermes/systemd/fleet-update.service
# Fleet update orchestrator — runs the Hermes fleet_update skill.
#
# Canonical source: tools/hermes/systemd/fleet-update.service in
# venturecrane/crane-console. Installed on mini by
# scripts/provision-hermes-fleet-update.sh (refuses to run on any other
# host). Paired with fleet-update.timer.
#
# See SKILL.md at tools/hermes/fleet_update/ and the plan at
# ~/.claude/plans/cuddly-riding-sifakis.md (#657).

[Unit]
Description=Crane fleet update orchestrator (Hermes)
Documentation=https://github.com/venturecrane/crane-console/tree/main/tools/hermes/fleet_update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=smdurgan
Group=smdurgan
WorkingDirectory=/srv/crane-console
Environment="PATH=/home/smdurgan/.local/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/etc/fleet-update/fleet-update.env

# Pin the repo checkout to origin/main on every run so edits in the
# canonical source take effect immediately and stale code is visible
# via the extra.source_sha field in ingested findings.
ExecStartPre=/usr/bin/git -C /srv/crane-console fetch --quiet
ExecStartPre=/usr/bin/git -C /srv/crane-console reset --quiet --hard origin/main

# Invocation form depends on the installed hermes-agent build. The
# provisioner substitutes ${HERMES_CMD} when installing if the env var
# is set; otherwise the default below is used. If the installed build
# does not support --skill / --non-interactive, override with:
#   HERMES_CMD='hermes skill run fleet_update' \
#   sudo bash scripts/provision-hermes-fleet-update.sh
ExecStart=/home/smdurgan/.local/bin/hermes chat --skill fleet_update --non-interactive

# One rolling log. Captured by systemd-journald anyway; the file is for
# easy Captain-tailing after a run.
StandardOutput=append:/var/log/fleet-update/run.log
StandardError=append:/var/log/fleet-update/run.log

# Hardening — we're running shell commands and SSH so a full sandbox
# is impractical, but prevent the service from gaining privileges or
# writing outside its expected surfaces.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/srv/crane-console /var/log/fleet-update /home/smdurgan
PrivateTmp=true

[Install]
WantedBy=multi-user.target
